# docker-tftp/supervisord.conf
[supervisord]
nodaemon=true

[program:in.tftpd]
command=/usr/sbin/in.tftpd --listen --create --user tftp --address 0.0.0.0:69 --secure /srv/tftp
